# deploy/linux/ghbrk.service
[Unit]
Description=ghbrk — privilege-separated git/gh broker
After=network.target

[Service]
Type=simple
User=ghbrk
Group=ghbrk-clients
ExecStart=/usr/local/bin/ghbrk daemon
Environment=GHBRK_POLICY=/etc/ghbrk/policy.yaml
Environment=GHBRK_AUDIT_LOG=/var/log/ghbrk/audit.log
Environment=GHBRK_SOCKET=/run/ghbrk/broker.sock
Environment=HOME=/run/ghbrk
Restart=on-failure
RestartSec=5s
# Hardening
NoNewPrivileges=true
AmbientCapabilities=CAP_SETUID CAP_SETGID
CapabilityBoundingSet=CAP_SETUID CAP_SETGID
ProtectSystem=strict
PrivateTmp=true
ReadWritePaths=/run/ghbrk /var/log/ghbrk
ProtectHome=no  # Required: user-impersonated child processes must write to repos under $HOME
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
